# crates/tessera_issuer/examples/issuer-serve.service
# systemd --user unit for the Tessera issuer local signing agent (Linux).
#
# Install (per user, no root):
#   mkdir -p ~/.config/systemd/user
#   cp issuer-serve.service ~/.config/systemd/user/
#   # edit the ExecStart flags below for your token module / key / cabinet origin
#   systemctl --user daemon-reload
#   systemctl --user enable --now issuer-serve.service
#
# The pairing token is written to $XDG_RUNTIME_DIR/tessera-issuer/token (0600);
# the cabinet reads it from there. The token PIN is entered via pinentry (a
# graphical prompt) or, absent that, the controlling terminal.

[Unit]
Description=Tessera issuer local signing agent
Documentation=https://github.com/TesseraLabs/tessera
After=default.target

[Service]
Type=simple
# Adjust the PKCS#11 module path, key label, algorithm and allowed origin(s).
ExecStart=%h/.local/bin/issuer serve \
    --daemon-token-file \
    --allow-origin https://cabinet.example \
    --module /usr/lib/x86_64-linux-gnu/opensc-pkcs11.so \
    --key tessera-ca \
    --algorithm ecdsa-p256
Restart=on-failure
# The agent listens only on 127.0.0.1 and needs no elevated privileges.
NoNewPrivileges=true

[Install]
WantedBy=default.target
